# apps/agent/examples/relaydock-agent.service
[Unit]
Description=RelayDock laptop agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/relaydock-agent run --config %h/.config/relaydock/agent.json
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
UMask=0077

[Install]
WantedBy=default.target
